# nmi.service.in
[Unit]
Description=Enable Open Power NMI service

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=@bindir@/pdbg -a stop
ExecStart=@bindir@/pdbg -a sreset
SyslogIdentifier=openpower-proc-nmi

