# op-cfam-reset.service.in
[Unit]
Description=CFAM reset
Wants=obmc-power-on@%i.target
After=obmc-power-on@%i.target
Wants=obmc-host-start-pre@0.target
Before=obmc-host-start-pre@0.target
Conflicts=obmc-host-stop@0.target
ConditionPathExists=!/run/openbmc/mpreboot@0

[Service]
RemainAfterExit=yes
Type=oneshot
ExecStart=/usr/bin/openpower-proc-control cfamReset

[Install]
# Will be installed in appropriate warm reboot targets via bb files
