# op-enter-mpreboot@.service.in
[Unit]
Description=Start memory preserving reboot host%i
Wants=obmc-host-stop-pre@%i.target
Before=obmc-host-stop-pre@%i.target
After=op-stop-instructions@%i.service
After=xyz.openbmc_project.Ipmi.Internal.SoftPowerOff.service
Conflicts=obmc-host-startmin@%i.target
ConditionPathExists=/run/openbmc/mpreboot@%i

[Service]
RemainAfterExit=yes
Type=oneshot
ExecStart=/usr/bin/openpower-proc-control enterMpipl

[Install]
WantedBy=obmc-host-stop@%i.target
