# op-stop-instructions@.service.in
[Unit]
Description=Stop instructions for host%i
Wants=obmc-host-stop-pre@%i.target
Before=obmc-host-stop-pre@%i.target
After=op-occ-disable@%i.service
After=xyz.openbmc_project.Ipmi.Internal.SoftPowerOff.service
Conflicts=obmc-host-startmin@%i.target
ConditionPathExists=/sys/class/fsi-master/fsi0/slave@00:00/cfam_id

[Service]
RemainAfterExit=yes
Type=oneshot
# TODO - Replace with official istep stopinstructions once available
ExecStart=/usr/bin/pdbg -a stop

[Install]
WantedBy=obmc-host-stop@%i.target
