# op-stop-instructions@.service.in
[Unit]
Description=Stop instructions for host%i
Wants=obmc-host-stop-pre@%i.target
After=obmc-host-stop-pre@%i.target
Wants=obmc-host-stopping@%i.target
Before=obmc-host-stopping@%i.target
After=op-occ-disable@%i.service
Conflicts=obmc-host-startmin@%i.target
ConditionPathExists=/sys/class/fsi-master/fsi0/slave@00:00/cfam_id
ConditionPathExists=!/run/openbmc/mpreboot@%i

[Service]
RemainAfterExit=yes
Type=oneshot
# TODO - Replace with official istep stopinstructions once available
ExecStart=/usr/bin/pdbg -a stop
ExecStart=/bin/rm -f /run/openbmc/host@%i-on

[Install]
WantedBy=obmc-host-stop@%i.target
