# pcie-poweroff@.service.in
[Unit]
Description=POWER9 PCIe Power-off Workaround
Before=op-power-stop@%i.service

[Service]
Type=oneshot
ExecStart=@sbindir@/openpower-proc-control cleanupPcie

[Install]
WantedBy=obmc-chassis-poweroff@%i.target