# set-spi-mux.service.in
[Unit]
Description=Set the SPI mux
ConditionPathExists=!/run/openbmc/chassis@0-on

[Service]
ExecStart=@bindir@/openpower-proc-control setSPIMux
SyslogIdentifier=openpower-proc-control
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
