# service_files/op-clock-data-logger@.service.in
[Unit]
Description=Host%i Clock debug data logger
Wants=obmc-host-started@%i.target
After=obmc-host-started@%i.target
After=openpower-update-bios-attr-table.service
Conflicts=obmc-host-stop@%i.target

[Service]
RemainAfterExit=yes
Type=simple
ExecStart=/usr/bin/openpower-clock-data-logger

[Install]
WantedBy=obmc-host-startmin@0.target
