# service_files/op-enter-mpreboot@.service.in
[Unit]
Description=Start memory preserving reboot host%i
Wants=obmc-host-stop-pre@%i.target
After=obmc-host-stop-pre@%i.target
Wants=obmc-host-stopping@%i.target
Before=obmc-host-stopping@%i.target
Before=op-stop-instructions@%i.service
After=phosphor-reset-sensor-states@%i.service
Before=phosphor-reboot-host@%i.service
After=openpower-update-bios-attr-table.service
Conflicts=obmc-host-startmin@%i.target

[Service]
@ENABLE_PHAL_TRUE@Environment="PDBG_DTB=/var/lib/phosphor-software-manager/hostfw/running/DEVTREE"
RemainAfterExit=yes
Type=oneshot
ExecStart=/usr/bin/openpower-proc-control enterMpReboot
ExecStart=/bin/mkdir -p /run/openbmc/
ExecStart=/bin/touch /run/openbmc/mpreboot@%i
ExecStart=/bin/sh -c "busctl set-property  xyz.openbmc_project.State.Host /xyz/openbmc_project/state/host%i xyz.openbmc_project.State.Host RestartCause s  xyz.openbmc_project.State.Host.RestartCause.HostCrash"

[Install]
WantedBy=obmc-host-diagnostic-mode@%i.target
