# service_files/op-enter-mpreboot@.service.in
[Unit]
Description=Start memory preserving reboot host%i
Wants=obmc-host-stop-pre@%i.target
After=obmc-host-stop-pre@%i.target
Wants=obmc-host-stopping@%i.target
Before=obmc-host-stopping@%i.target
Before=op-stop-instructions@%i.service
After=phosphor-reset-sensor-states@%i.service
Before=phosphor-reboot-host@%i.service
Before=clear_hostdumps_poweroff.service
After=openpower-update-bios-attr-table.service
Conflicts=obmc-host-startmin@%i.target

[Service]
@ENABLE_PHAL_TRUE@Environment="PDBG_DTB=@CEC_DEVTREE_RW_PATH@"
RemainAfterExit=yes
Type=oneshot
ExecStart=/bin/sh -c \
  "busctl set-property xyz.openbmc_project.BIOSConfigManager \
   /xyz/openbmc_project/bios_config/manager xyz.openbmc_project.BIOSConfig.Manager \
   PendingAttributes a{s\(sv\)} 1 pvm_sys_dump_active \
   xyz.openbmc_project.BIOSConfig.Manager.AttributeType.Enumeration s Enabled"
ExecStart=/bin/sh -c \
  "busctl call  xyz.openbmc_project.Dump.Manager /xyz/openbmc_project/dump/bmc \
   xyz.openbmc_project.Dump.Create CreateDump a{sv} 1 \
   com.ibm.Dump.Create.CreateParameters.DumpType s com.ibm.Dump.Create.DumpType.BMC" || true
ExecStart=/usr/bin/openpower-proc-control enterMpReboot
ExecStart=/bin/mkdir -p /run/openbmc/
ExecStart=/bin/touch /run/openbmc/mpreboot@%i
ExecStart=/bin/sh -c "busctl set-property  xyz.openbmc_project.State.Host /xyz/openbmc_project/state/host%i xyz.openbmc_project.State.Host RestartCause s  xyz.openbmc_project.State.Host.RestartCause.HostCrash"

[Install]
WantedBy=obmc-host-diagnostic-mode@%i.target
