# service_files/op-reset-host-check@.service.in
[Unit]
Description=Check Host%i status on BMC reset via CFAM reg
Wants=mapper-wait@-xyz-openbmc_project-state-host%i.service
After=mapper-wait@-xyz-openbmc_project-state-host%i.service
Wants=mapper-wait@-xyz-openbmc_project-logging.service
After=mapper-wait@-xyz-openbmc_project-logging.service
Wants=obmc-host-reset-running@%i.target
Before=obmc-host-reset-running@%i.target
Wants=op-reset-chassis-on@%i.service
After=op-reset-chassis-on@%i.service
After=openpower-update-bios-attr-table.service
Conflicts=obmc-host-stop@%i.target
ConditionPathExists=/run/openbmc/chassis@%i-on
ConditionPathExists=!/run/openbmc/host@%i-on

[Service]
RemainAfterExit=yes
Type=oneshot
ExecStart=/usr/bin/openpower-proc-control checkHostRunning

[Install]
WantedBy=obmc-host-reset@%i.target
