# service_files/op-reset-host-clear.service.in
[Unit]
Description=Clear host running status register
Wants=obmc-host-stop-pre@0.target
Before=obmc-host-stop-pre@0.target
Conflicts=obmc-host-startmin@0.target

[Service]
RemainAfterExit=yes
Type=oneshot
ExecStart=/usr/bin/putcfam pu 2983 -pall 0

[Install]
WantedBy=obmc-host-stop@0.target
