# service_files/op-stop-instructions@.service.in
[Unit]
Description=Stop instructions for host%i
After=obmc-host-stop-pre@%i.target
Before=obmc-host-stopping@%i.target
After=op-occ-disable@%i.service
After=openpower-update-bios-attr-table.service
Conflicts=obmc-host-startmin@%i.target
ConditionPathExists=/sys/class/fsi-master/fsi0/slave@00:00/cfam_id
ConditionPathExists=!/run/openbmc/mpreboot@%i

[Service]
RemainAfterExit=yes
Type=oneshot
TimeoutStartSec=20
ExecStart=/usr/bin/openpower-proc-control threadStopAll

[Install]
WantedBy=obmc-host-stop@%i.target
