# service_files/pcie-poweroff@.service.in
[Unit]
Description=POWER9 PCIe Power-off Workaround
Before=phosphor-power-stop@%i.service
After=obmc-power-stop-pre@%i.target

[Service]
Type=oneshot
ExecStart=@bindir@/openpower-proc-control cleanupPcie

[Install]
WantedBy=obmc-chassis-poweroff@%i.target
