# service_files/set-spi-mux.service.in
[Unit]
Description=Set the SPI mux
ConditionPathExists=!/run/openbmc/chassis@0-on
After=op-reset-chassis-running@0.service
After=op-wait-power-off@0.service

[Service]
ExecStartPre=@bindir@/openpower-proc-control cfamReset
ExecStartPre=@bindir@/openpower-proc-control scanFSI
ExecStart=@bindir@/openpower-proc-control setSPIMux
ExecStartPost=@bindir@/openpower-proc-control scanFSI
SyslogIdentifier=openpower-proc-control
Type=oneshot

[Install]
WantedBy=multi-user.target
WantedBy=obmc-chassis-poweroff@0.target
